# bridge/com.tikita.bambu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.tikita.bambu</string>
  <!-- caffeinate -s: köprü çalışırken Mac (fişe takılıyken) uykuya dalmaz -->
  <!-- /usr/bin/python = macOS 11 sistem Python 2.7.16. python3 kullaniyorsan bunu /usr/bin/python3 yap. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-s</string>
    <string>/usr/bin/python</string>
    <string>/Users/KULLANICI_ADIN/tikita/bambu_kopru.py</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/tikita-bambu.log</string>
  <key>StandardErrorPath</key><string>/tmp/tikita-bambu.log</string>
</dict>
</plist>
